# debian/postinst.in
#!/bin/sh
# NightWatcher2 post-install: PATH/systemd/udev symlinks, /etc config seed,
# debconf-driven database + API configuration, USB (dialout) permissions, and
# systemd enablement.
set -e

PREFIX="@CMAKE_INSTALL_PREFIX@"
SYMLINK_BINDIR="@NW_SYMLINK_BINDIR@"
CONF=/etc/nightwatcher/nightwatcher.conf
ENVF=/etc/nightwatcher/nightwatcher.env

HAVE_DEBCONF=false
if [ -f /usr/share/debconf/confmodule ]; then
  . /usr/share/debconf/confmodule
  HAVE_DEBCONF=true
fi

# Set "KEY = VALUE" under [SECTION] in an INI file, section-aware: [database]
# and [api] both have a `port` key, so a blind sed would rewrite the wrong one.
set_ini() {  # set_ini FILE SECTION KEY VALUE
  awk -v sec="$2" -v key="$3" -v val="$4" '
    /^[[:space:]]*\[/ { insec = ($0 ~ ("^[[:space:]]*\\[" sec "\\]")) }
    insec && $0 ~ ("^[[:space:]]*" key "[[:space:]]*=") { print key " = " val; next }
    { print }
  ' "$1" > "$1.nwtmp" && mv "$1.nwtmp" "$1"
}

case "$1" in
  configure)
    # 1) CLI symlinks on PATH (the daemon is started by systemd via absolute path).
    mkdir -p "$SYMLINK_BINDIR"
    for tool in sqmctl nwdb nwexport-auth; do
      [ -x "$PREFIX/bin/$tool" ] && ln -sf "$PREFIX/bin/$tool" "$SYMLINK_BINDIR/$tool"
    done

    # 2) systemd unit + udev rule into their searched paths.
    ln -sf "$PREFIX/lib/systemd/system/nightwatcherd.service" /etc/systemd/system/nightwatcherd.service
    ln -sf "$PREFIX/lib/udev/rules.d/99-nightwatcher-sqm.rules" /etc/udev/rules.d/99-nightwatcher-sqm.rules

    # 3) Seed /etc config, never clobbering an existing one.
    mkdir -p /etc/nightwatcher
    if [ ! -e "$CONF" ]; then
      cp "$PREFIX/share/nightwatcher.conf.example" "$CONF"
    fi

    # 4) Apply the debconf answers: DB connection + API settings into the config,
    #    the password into the 0600 env file, and (optionally) create the DB.
    DB_READY=false
    START_NOW=false
    if [ "$HAVE_DEBCONF" = true ]; then
      db_get nightwatcher/db_password; DB_PASS="$RET"
      db_get nightwatcher/db_host;     DB_HOST="$RET"
      db_get nightwatcher/db_create;   DB_CREATE="$RET"
      db_get nightwatcher/api_bind;    API_BIND="$RET"
      db_get nightwatcher/api_port;    API_PORT="$RET"
      db_get nightwatcher/api_tls;     API_TLS="$RET"
      db_get nightwatcher/start_now;   START_NOW="$RET"
      # Release the debconf protocol before any plain output below.
      db_stop

      if [ -f "$CONF" ]; then
        [ -n "$DB_HOST" ]  && set_ini "$CONF" database host "$DB_HOST"
        [ -n "$API_BIND" ] && set_ini "$CONF" api bind "$API_BIND"
        [ -n "$API_PORT" ] && set_ini "$CONF" api port "$API_PORT"
        case "$API_TLS" in
          true)  set_ini "$CONF" api tls on  ;;
          false) set_ini "$CONF" api tls off ;;
        esac
      fi

      # Secret -> env file (0600, read by the systemd unit as NW_DB_PASSWORD).
      if [ -n "$DB_PASS" ]; then
        umask 077
        tmp=$(mktemp "${ENVF}.XXXXXX")
        if [ -f "$ENVF" ]; then
          grep -v '^[[:space:]]*NW_DB_PASSWORD=' "$ENVF" >> "$tmp" 2>/dev/null || true
        fi
        printf 'NW_DB_PASSWORD=%s\n' "$DB_PASS" >> "$tmp"
        mv "$tmp" "$ENVF"
        chmod 0600 "$ENVF"
      fi

      # Optionally create the database + user via the local root socket. The
      # daemon builds the tables and seeds the admin account on first start.
      is_local=false
      case "$DB_HOST" in
        ""|localhost|127.0.0.1|::1|localhost:*|127.0.0.1:*) is_local=true ;;
      esac
      if [ "$DB_CREATE" = true ] && [ "$is_local" = true ] && [ -n "$DB_PASS" ]; then
        MYSQL=
        for c in mariadb mysql; do
          command -v "$c" >/dev/null 2>&1 && { MYSQL="$c"; break; }
        done
        if [ -n "$MYSQL" ] && "$MYSQL" -e 'SELECT 1' >/dev/null 2>&1; then
          esc=$(printf '%s' "$DB_PASS" | sed "s/'/''/g")
          if "$MYSQL" >/dev/null 2>&1 <<SQL
CREATE DATABASE IF NOT EXISTS \`nightwatcher\` CHARACTER SET utf8mb4 COLLATE utf8mb4_unicode_ci;
CREATE USER IF NOT EXISTS 'nightwatcher'@'localhost' IDENTIFIED BY '$esc';
CREATE USER IF NOT EXISTS 'nightwatcher'@'127.0.0.1' IDENTIFIED BY '$esc';
ALTER USER 'nightwatcher'@'localhost' IDENTIFIED BY '$esc';
ALTER USER 'nightwatcher'@'127.0.0.1' IDENTIFIED BY '$esc';
GRANT ALL PRIVILEGES ON \`nightwatcher\`.* TO 'nightwatcher'@'localhost';
GRANT ALL PRIVILEGES ON \`nightwatcher\`.* TO 'nightwatcher'@'127.0.0.1';
FLUSH PRIVILEGES;
SQL
          then
            DB_READY=true
            echo "nightwatcher: database 'nightwatcher' and user are ready."
          else
            echo "nightwatcher: automatic database setup failed; run 'sudo mariadb < $PREFIX/sql/setup.sql' yourself."
          fi
        else
          echo "nightwatcher: no reachable local MariaDB/MySQL root socket; run 'sudo mariadb < $PREFIX/sql/setup.sql' yourself."
        fi
      elif [ "$is_local" != true ] && [ -n "$DB_PASS" ]; then
        DB_READY=true    # remote DB: assume the operator manages it
      elif [ -n "$DB_PASS" ]; then
        DB_READY=true    # local, creation declined: operator takes responsibility
      fi
    fi

    # 5) USB SQM-LU access: apply the udev rule now, and give the installing
    #    user serial access so they can run sqmctl interactively.
    if command -v udevadm >/dev/null 2>&1; then
      udevadm control --reload-rules >/dev/null 2>&1 || true
      udevadm trigger --subsystem-match=tty >/dev/null 2>&1 || true
    fi
    if [ -n "${SUDO_USER:-}" ] && [ "${SUDO_USER}" != root ]; then
      if ! id -nG "${SUDO_USER}" 2>/dev/null | tr ' ' '\n' | grep -qx dialout; then
        usermod -aG dialout "${SUDO_USER}" 2>/dev/null || true
        echo "nightwatcher: added ${SUDO_USER} to group 'dialout' (log out/in for it to take effect)."
      fi
    fi

    # 6) systemd: load + enable the unit; start it only when the DB looks ready.
    if command -v systemctl >/dev/null 2>&1; then
      systemctl daemon-reload >/dev/null 2>&1 || true
      systemctl enable nightwatcherd.service >/dev/null 2>&1 || true
      if [ "${START_NOW:-false}" = true ] && [ "$DB_READY" = true ]; then
        if systemctl restart nightwatcherd.service >/dev/null 2>&1; then
          echo "nightwatcher: service started (web UI on ${API_BIND:-127.0.0.1}:${API_PORT:-8080})."
        else
          echo "nightwatcher: service failed to start; check 'journalctl -u nightwatcherd'."
        fi
      fi
    fi

    if [ "${START_NOW:-false}" != true ] || [ "$DB_READY" != true ]; then
      echo "nightwatcher: when the database is ready -> sudo systemctl start nightwatcherd"
    fi
    ;;
esac

exit 0
