# debian/prerm.in
#!/bin/sh
# NightWatcher2 pre-removal: stop and disable the service before files go away.
set -e
case "$1" in
  remove|deconfigure)
    if command -v systemctl >/dev/null 2>&1; then
      systemctl disable --now nightwatcherd.service >/dev/null 2>&1 || true
    fi
    ;;
esac
exit 0
